# scripts/ritual-aurora-clipboard-trim.timer
[Unit]
Description=Trim diário do histórico do clipboard manager (mantém últimas 200 entradas)

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=15min

[Install]
WantedBy=timers.target
